# deploy/vllm-qwen-32b-fp8-2.service
# vLLM 实例 2 - Qwen3-VL-32B-Instruct-FP8 (GPU 2,3)
# 安装: 见 vllm-qwen-32b-fp8-1.service

[Unit]
Description=vLLM Qwen-32B-FP8 Instance 2 (GPU 2,3)
After=network.target

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/home/azureuser/workspace/tag

Environment=CUDA_VISIBLE_DEVICES=2,3
Environment=OMP_NUM_THREADS=1

ExecStart=/anaconda/envs/qwen-vl/bin/vllm serve /home/azureuser/models/Qwen3-VL-32B-Instruct-FP8 \
  --tensor-parallel-size 2 \
  --port 8004 \
  --host 0.0.0.0 \
  --max-model-len 131072

Restart=on-failure
RestartSec=30
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
